# client/build/pkg-scripts/postinstall
#!/bin/bash

APP_PATH="/Applications/Proxyness.app"
HELPER_DIR="/Library/PrivilegedHelperTools"
HELPER_BIN="$HELPER_DIR/proxyness-helper"
PLIST_PATH="/Library/LaunchDaemons/com.proxyness.helper.plist"

# --- Migration from pre-rebrand names ---
OLD_LABEL="com.smurov.proxy.helper"
OLD_PLIST="/Library/LaunchDaemons/${OLD_LABEL}.plist"
OLD_BIN="$HELPER_DIR/smurov-helper"

launchctl bootout "system/${OLD_LABEL}" 2>/dev/null || true
rm -f "$OLD_PLIST" "$OLD_BIN"
rm -f /var/run/smurov-helper.sock /var/run/proxyness-helper.sock
# --- End migration ---

# Clear quarantine and fix permissions on bundled binaries
xattr -cr "$APP_PATH" 2>/dev/null || true
chmod +x "$APP_PATH/Contents/Resources/resources"/daemon-* 2>/dev/null || true
chmod +x "$APP_PATH/Contents/Resources/resources"/helper-* 2>/dev/null || true

# Install helper binary
mkdir -p "$HELPER_DIR"
HELPER_SRC=$(find "$APP_PATH/Contents/Resources/resources" -name "helper-darwin-*" | head -1)
if [ -n "$HELPER_SRC" ]; then
    # Unload current version if running
    launchctl bootout system/com.proxyness.helper 2>/dev/null || true

    cp "$HELPER_SRC" "$HELPER_BIN"
    chmod 755 "$HELPER_BIN"
    chown root:wheel "$HELPER_BIN"

    # Install launchd plist
    cat > "$PLIST_PATH" << 'PLIST'
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.proxyness.helper</string>
    <key>Program</key>
    <string>/Library/PrivilegedHelperTools/proxyness-helper</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
PLIST

    chmod 644 "$PLIST_PATH"
    chown root:wheel "$PLIST_PATH"

    # Start the daemon
    launchctl bootstrap system "$PLIST_PATH"
fi

exit 0
